# deploy/systemd/xicemc-maintenance.service
[Unit]
Description=XiceMCServer Daily Backup and Update
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
Group=root
WorkingDirectory=/opt/xicemc/repo
EnvironmentFile=-/opt/xicemc/secrets/xicemc.env
ExecStart=/opt/xicemc/repo/scripts/daily-maintenance.sh
