# deploy/systemd/xicemc-whitelist.service
[Unit]
Description=XiceMCServer Web Portal
After=network-online.target xicemc.service
Wants=network-online.target

[Service]
Type=simple
User=minecraft
Group=minecraft
WorkingDirectory=/opt/xicemc/repo
EnvironmentFile=/opt/xicemc/secrets/whitelist-web.env
EnvironmentFile=-/opt/xicemc/secrets/xicemc.env
ExecStart=/opt/xicemc/bin/xicemc-web-go
Restart=on-failure
RestartSec=10
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
